# newsroom/com.la-aurora.watchdog.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
  <dict>
    <key>Label</key>
    <string>com.la-aurora.watchdog</string>
    <key>ProgramArguments</key>
    <array>
      <string>/usr/bin/python3</string>
      <string>/Users/c/Library/LaAurora/newsroom/watchdog_monitor.py</string>
    </array>
    <key>RunAtLoad</key>
    <true/>
    <key>StartInterval</key>
    <integer>120</integer>
    <key>KeepAlive</key>
    <false/>
    <key>StandardOutPath</key>
    <string>/Users/c/Library/LaAurora/newsroom/logs/watchdog.out.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/c/Library/LaAurora/newsroom/logs/watchdog.err.log</string>
  </dict>
</plist>
